# content/service/redis/run
#!/bin/sh

# Run redis
exec 2>&1
exec redis-server --port 61801
