# content/service/aria2/run
#!/bin/sh

# Configure aria2
touch /workdir/aria2.session
/workdir/tracker.sh /workdir/aria2.conf

# Run aria2
exec 2>&1
exec aria2c --conf-path="/workdir/aria2.conf" --enable-rpc